# debian/acr.service
[Unit]
Description=AudioControl Rust (ACR) - Audio Player Management Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/acr -c /etc/acr/acr.json
Restart=on-failure
RestartSec=5
User=acr
Group=acr
WorkingDirectory=/etc/acr
RuntimeDirectory=acr
RuntimeDirectoryMode=0755
PIDFile=/var/run/acr/acr.pid
TimeoutStopSec=10s
KillSignal=SIGTERM
KillMode=process

[Install]
WantedBy=multi-user.target
