# debian/audiocontrol.service
[Unit]
Description=HiFiBerry AudioControl - Audio Player Management Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/audiocontrol -c /etc/audiocontrol/audiocontrol.json
Restart=on-failure
RestartSec=5
User=audiocontrol
Group=audiocontrol
WorkingDirectory=/etc/audiocontrol
RuntimeDirectory=audiocontrol
RuntimeDirectoryMode=0755
PIDFile=/var/run/audiocontrol/audiocontrol.pid

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
