# debian/postinst
#!/bin/bash
set -e

# Ensure the script runs with proper error handling (bash specific)
if [ -n "$BASH_VERSION" ]; then
    set -euo pipefail
fi

# Function to safely copy default config if none exists
setup_config() {
    # Check if we need to copy the sample config to the actual config file
    if [ ! -f /etc/acr/acr.json ]; then
        echo "No configuration file found, copying from sample..."
        if [ -f /etc/acr/acr.json.sample ]; then
            cp /etc/acr/acr.json.sample /etc/acr/acr.json
            # Set proper ownership
            chown acr:acr /etc/acr/acr.json
            echo "Sample config copied to acr.json successfully"
        else
            # Try to find the sample config elsewhere if it's not in /etc/acr
            if [ -f /usr/share/acr/acr-default.json ]; then
                echo "Using default configuration from /usr/share/acr..."
                cp /usr/share/acr/acr-default.json /etc/acr/acr.json
                chown acr:acr /etc/acr/acr.json
                echo "Default config copied to acr.json successfully"
            else
                echo "ERROR: No sample configuration file found"
                echo "Cannot continue without a valid configuration file"
                exit 1
            fi
        fi
    else
        echo "Existing configuration file found, keeping it."
    fi
}

# Function to create user, group and directory
setup_user_and_dirs() {
    # Handle legacy installations at /acr
    if [ -d /acr ]; then
        echo "Found legacy directory at /acr, migrating to /etc/acr..."
        
        # Copy config files from /acr to /etc/acr if they exist
        if [ -f /acr/acr.json ] && [ ! -f /etc/acr/acr.json ]; then
            echo "Migrating configuration from /acr/acr.json to /etc/acr/acr.json..."
            cp /acr/acr.json /etc/acr/acr.json
        fi
        
        # No need to remove /acr as it might be used by other components
    fi

    # Create acr group if it doesn't exist
    if ! getent group acr > /dev/null; then
        echo "Creating acr group..."
        groupadd --system acr
    fi
    
    # Create acr user if it doesn't exist
    if ! getent passwd acr > /dev/null; then
        echo "Creating acr user..."
        useradd --system --gid acr --shell /usr/sbin/nologin --home-dir /etc/acr acr
    fi    
    
    # Create directory structure
    echo "Setting up directory structure..."
    mkdir -p /etc/acr/cache/attributes
    mkdir -p /etc/acr/cache/images
    mkdir -p /var/acr
    
    # Fix permissions for all directories
    echo "Setting proper permissions..."
    chown -R acr:acr /etc/acr
    chmod 755 /etc/acr
    
    chown -R acr:acr /var/acr
    chmod 755 /var/acr
}

case "$1" in
    configure)
        setup_user_and_dirs
        setup_config
        # Enable and start the service
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload >/dev/null 2>&1 || true
            systemctl enable acr.service >/dev/null 2>&1 || true
            systemctl restart acr.service >/dev/null 2>&1 || true
        fi
        ;;
esac

exit 0
